# .config/systemd/user/polybar.service
[Unit]
Description=Polybar system status bar

[Service]
ExecStart=/usr/bin/polybar main
Type=process
KillMode=process

