# .config/systemd/user/xsession.service
[Unit]
Description=Xsession dummy service
PartOf=xsession.target


[Service]
ExecStart=/usr/bin/echo Xsession started
Type=oneshot
